# meta-zvezda/recipes-phosphor/manufacturer-settings/files/boardinfo.timer
[Unit]
Description=Timer runs boardinfo

[Timer]
OnBootSec=2min
Unit=boardinfo.service

[Install]
WantedBy=multi-user.target
