# meta-zvezda/recipes-phosphor/manufacturer-settings/files/boardinfo.service
[Unit]
Description=Boardinfo setter

[Service]
ExecStart=/usr/bin/boardinfo
Type=one-shot
