# meta-zvezda/meta-iridium64/recipes-phosphor/extra-configs/files/set-gpios.service
[Unit]
Description=Set configuration gpios
Before=xyz.openbmc_project.Inventory.Manager.service

[Service]
ExecStart=/usr/bin/set-gpios
Type=simple

[Install]
WantedBy=multi-user.target
