# meta-zvezda/meta-iridium64/recipes-phosphor/extra-configs/files/set-gpios.timer
[Unit]
Description=Set configuration gpios

[Timer]
OnBootSec=0sec
Unit=set-gpios.service

[Install]
WantedBy=multi-user.target
