# meta-zvezda/meta-iridium64/recipes-phosphor/extra-configs/files/set-gpios.service
[Unit]
Description=Set configuration gpios

[Service]
ExecStart=/usr/bin/set-gpios
Type=one-shot